# rootfs/etc/systemd/system/munge.service.d/override.conf
[Service]
# Self-correct ownership on our munge keys and config.
ExecStartPre=+/usr/bin/chown munge:munge -R /etc/munge
ExecStartPre=+/usr/bin/chmod 0600 /etc/munge/munge.key
